# builder_global/global/profile.d/profile.jenkins.ant
<?xml version="1.0" encoding="UTF-8"?>
<project name="openscada.build.profile">
	<echo message="Profile: Jenkins CI" />
	
	<property environment="env"/>
	
	<property name="eclipse.home.dir" location="${env.ECLIPSE_431_HOME}"/>
	<property name="openscada.mirror.base" value="${env.OPENSCADA_MIRROR_BASE}"/>
	
	<property name="build.root" location="../.."/>
	<property name="jobSuffix" value="1.2"/>
	
	<property name="output.dir" location="${env.WORKSPACE}/output"/>
	<property name="input.dir" location="${env.WORKSPACE}/input"/>
	<property name="jenkins.download.base" value="${env.JENKINS_URL}"/>
	
	<!-- load properties last, since they don't override in ant -->
	<property file="${build.root}/global.properties"/>
	<property file="${build.root}/global/global.properties"/>
	
	<import file="${build.root}/global/tasks.ant"/>
	
	<!-- ================================= 
          target: setup
         ================================= -->
    <target name="setup" unless="skip.setup">
    	
    	<delete dir="${input.dir}" failonerror="false"/>
		<mkdir dir="${input.dir}"/>
    	
    	<getP2Repo required="${require.utils}" jobname="org.eclipse.scada.utils-master" reponame="org.eclipse.scada.utils.p2-incubation" version="0.1.0" qualifier="-SNAPSHOT" suffix="" server="https://hudson.eclipse.org/scada"/>
    	<getP2Repo required="${require.base}" jobname="org.eclipse.scada.base-master" reponame="org.eclipse.scada.base.p2-incubation" version="0.1.0" qualifier="-SNAPSHOT" suffix="" server="https://hudson.eclipse.org/scada"/>
    	<getP2Repo required="${require.infinity}" jobname="org.eclipse.scada.chart-master" reponame="org.eclipse.scada.chart.p2-incubation" version="0.1.0" qualifier="-SNAPSHOT" suffix="" server="https://hudson.eclipse.org/scada"/>
    	<getP2Repo required="${require.protocols}" jobname="org.eclipse.scada.protocols-master" reponame="org.eclipse.scada.protocols.p2-incubation" version="0.1.0" qualifier="-SNAPSHOT" suffix="" server="https://hudson.eclipse.org/scada"/>
    	<getP2Repo required="${require.core}" jobname="org.eclipse.scada.core-master" reponame="org.eclipse.scada.core.p2-incubation" version="0.1.0" qualifier="-SNAPSHOT" suffix="" server="https://hudson.eclipse.org/scada"/>
    	<getP2Repo required="${require.hmi}" jobname="org.eclipse.scada.hmi-master" reponame="org.eclipse.scada.hmi.p2-incubation" version="0.1.0" qualifier="-SNAPSHOT" suffix="" server="https://hudson.eclipse.org/scada"/>
    	
    	<getP2Repo required="${require.external}" jobname="org.eclipse.scada.external-master" reponame="org.eclipse.scada.external.p2-incubation" version="0.1.0" qualifier="-SNAPSHOT" suffix="" server="https://hudson.eclipse.org/scada"/>
    	<getP2Repo required="${require.external}" jobname="org.openscada.external-${jobSuffix}" reponame="org.openscada.external.sdk" version="${external.version}"/>

		<getP2Repo required="${require.utgard}" jobname="org.openscada.utgard-${jobSuffix}" reponame="org.openscada.utgard.sdk" version="${utgard.version}"/>
    	<getP2Repo required="${require.atlantis}" jobname="org.openscada.atlantis-${jobSuffix}" reponame="org.openscada.atlantis.sdk" version="${atlantis.version}"/>
    	<getP2Repo required="${require.deploy}" jobname="org.openscada.deploy-${jobSuffix}" reponame="org.openscada.deploy.repo" version="${deploy.version}"/>
    	
    	<getP2Repo required="${require.jinterop}" jobname="org.openscada.jinterop-${jobSuffix}" reponame="org.openscada.jinterop.sdk" version="${jinterop.version}"/>
    	<getP2Repo required="${require.docdocdoc}" jobname="org.openscada.docdocdoc-${jobSuffix}" reponame="org.openscada.docdocdoc.sdk" version="${docdocdoc.version}"/>
    	<getP2Repo required="${require.contrib}" jobname="org.openscada.contrib-${jobSuffix}" reponame="org.openscada.contrib.sdk" version="${contrib.version}"/>
    	
    	<getProduct required="${require.p2director}" jobname="org.openscada.deploy.p2director-${jobSuffix}" productName="p2director" version="${deploy.version}"/>
    	
    	<downloadP2Repo required="${require.external.rap}" uri="${openscada.mirror.base}/rap-2.1.0-R-20130611-2139.zip"/>
    	<downloadP2Repo required="${require.external.wtp}" uri="${openscada.mirror.base}/wtp4x-repo-R-3.4.2-20130208151217.zip"/>
    	<downloadP2Repo required="${require.external.gef}" uri="${openscada.mirror.base}/GEF-Update-3.9.0.zip"/>
    	<downloadP2Repo required="${require.external.emf}" uri="${openscada.mirror.base}/emf-xsd-Update-2.9.0.zip"/>
		<downloadP2Repo required="${require.external.rcp}" uri="${openscada.mirror.base}/org.eclipse.rcp-4.3.zip"/>
		<downloadP2Repo required="${require.external.platform}" uri="${openscada.mirror.base}/org.eclipse.platform-4.3.zip"/>
    	<downloadP2Repo required="${require.external.ecf}" uri="${openscada.mirror.base}/org.eclipse.ecf.sdk_3.5.6.v20120610-1946.zip"/>
    	
    	<downloadP2Repo required="${require.external.equinox}" uri="${openscada.mirror.base}/equinox-SDK-Kepler.zip"/>
    	<downloadP2Repo required="${require.external.gemini}" uri="${openscada.mirror.base}/gemini-dbaccess-1.1.0.RELEASE-updates.zip"/>
    	
    	<downloadEclipseArchive required="${require.external.babel}" uri="${openscada.mirror.base}/BabelLanguagePack-eclipse-de_4.2.0.v20121120043402.zip"/>
    	<downloadEclipseArchive required="${require.external.babel}" uri="${openscada.mirror.base}/BabelLanguagePack-rt.equinox-de_4.2.0.v20121120043402.zip"/>

    </target>
	
	<!-- = = = = = = = = = = = = = = = = =
          macrodef: downloadEclipseArchive          
         = = = = = = = = = = = = = = = = = -->
    <macrodef name="downloadEclipseArchive">
    	<attribute name="required" />
    	<attribute name="uri"/>
        <sequential>
        	<if>
        		<equals arg1="@{required}" arg2="true"/>
        		<then>
        			<mkdir dir="${input.dir}/eclipse"/>
        			<echo message="Importing external Eclipse Archive (non-P2 Repository) from: @{uri}"/>
			    	<get dest="${input.dir}/eclipse" src="@{uri}" usetimestamp="true"/>
        		</then>
        	</if>
        </sequential>
    </macrodef>
	
	<!-- = = = = = = = = = = = = = = = = =
          macrodef: downloadP2Repo          
         = = = = = = = = = = = = = = = = = -->
    <macrodef name="downloadP2Repo">
    	<attribute name="required" />
    	<attribute name="uri"/>
        <sequential>
        	<if>
        		<equals arg1="@{required}" arg2="true"/>
        		<then>
        			<mkdir dir="${input.dir}/repos"/>
        			<echo message="Importing external P2 Repository from: @{uri}"/>
			    	<get dest="${input.dir}/repos" src="@{uri}" usetimestamp="true"/>
        		</then>
        	</if>
        </sequential>
    </macrodef>

	
	<!-- = = = = = = = = = = = = = = = = =
          macrodef: getP2Repo          
         = = = = = = = = = = = = = = = = = -->
    <macrodef name="getP2Repo">
    	<attribute name="server" default="${jenkins.download.base}" />
        <attribute name="repoName" />
    	<attribute name="version" />
    	<attribute name="qualifier" default=""/>
    	<attribute name="jobName" default="@{repoName}-@{version}" />
    	<attribute name="required" />
    	<attribute name="suffix" default="-p2repo"/>
        <sequential>
        	<if>
        		<equals arg1="@{required}" arg2="true"/>
        		<then>
        			<echo message="Importing P2 Repository from: @{server}/job/@{jobName}/lastSuccessfulBuild/artifact/output/@{repoName}-@{version}@{qualifier}@{suffix}.zip"/>
        			<mkdir dir="${input.dir}/repos"/>
			    	<get dest="${input.dir}/repos" src="@{server}/job/@{jobName}/lastSuccessfulBuild/artifact/output/@{repoName}-@{version}@{qualifier}@{suffix}.zip" usetimestamp="true"/>
        		</then>
        	</if>
        </sequential>
    </macrodef>
	
	<!-- = = = = = = = = = = = = = = = = =
          macrodef: getProduct          
         = = = = = = = = = = = = = = = = = -->
    <macrodef name="getProduct">
        <attribute name="productName" />
    	<attribute name="version" />
    	<attribute name="required" />
    	<attribute name="jobName" />
        <sequential>
        	<if>
        		<equals arg1="@{required}" arg2="true"/>
        		<then>
        			<echo message="Importing product @{productName} from: ${jenkins.download.base}"/>
                	<foreach list="${configs}" delimiter="&amp;" param="config" trim="true" target="targetGetProduct">
                		<param name="jobName" value="@{jobName}"/>
                		<param name="productName" value="@{productName}"/>
                		<param name="version" value="@{version}"/>
                	</foreach>
        		</then>
        	</if>
        </sequential>
    </macrodef>
	

	<!-- - - - - - - - - - - - - - - - - - 
          target: targetGetProduct                      
         - - - - - - - - - - - - - - - - - -->
        <target name="targetGetProduct">
    	<propertyregex property="configCleaned1" input="${config}" regexp=" +" replace=""/>
    	<propertyregex property="configCleaned" input="${configCleaned1}" regexp="," replace="."/>
    	<echo message="Processing config: ${configCleaned}"/>
    	<mkdir dir="${input.dir}/product"/>
		<get
			dest="${input.dir}/product"
			src="${jenkins.download.base}/job/${jobName}/lastSuccessfulBuild/artifact/output/${productName}/${productName}-${configCleaned}.zip"
			usetimestamp="true"/>
    </target>

</project>
